# supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:php-fpm]
command = /usr/sbin/php-fpm -F
autostart=true
autorestart=true

[program:td-agent]
command = /usr/sbin/td-agent -q
autostart=true
autorestart=true

[program:crond]
command = /usr/sbin/crond -n
autostart=true
autorestart=true
